# postgresql/postgresql.ini
[program:postgresql]
autorestart=true
autostart=true
command="/usr/bin/postgres" -D "%(ENV_PG_DATA)s" -p "%(ENV_PG_PORT)s"
priority=0
redirect_stderr=true
stdout_logfile=/var/log/pgsql.log
stdout_logfile_backups=5
stdout_logfile_maxbytes=50MB
user=postgres
