# opennms/opennms.ini
[program:opennms]
autorestart=true
autostart=true
command=/opt/opennms/bin/opennms -f start
priority=0
redirect_stderr=true
stdout_logfile=/data/opennms/logs/stdout.log
stdout_logfile_backups=5
stdout_logfile_maxbytes=50MB
